# Future/binlog_purge.timer
[Unit]
Description=Run binlog backup & purge script daily

[Timer]
OnCalendar=*-*-* 10:00:00
Persistent=true

[Install]
WantedBy=timers.target
